# scripts/picalender.service
[Unit]
Description=PiCalendar Clock Kiosk (pygame KMSDRM)
After=multi-user.target network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=zabaglione
Group=zabaglione
WorkingDirectory=/home/zabaglione/picalender
Environment="SDL_VIDEODRIVER=kmsdrm"
Environment="SDL_FBDEV=/dev/fb0"
Environment="SDL_KMSDRM_DEVICE_INDEX=0"
Environment="SDL_RENDER_DRIVER=software"
Environment="DISPLAY=:0"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/zabaglione/picalender:/home/zabaglione/picalender/src:/home/zabaglione/picalender/src/renderers"

# 起動前の待機時間（ネットワークとディスプレイの準備を短縮）
ExecStartPre=/bin/sleep 3

# メインアプリケーション起動（起動スクリプト経由）
ExecStart=/home/zabaglione/picalender/scripts/start_service.sh

# 再起動設定
Restart=on-failure
RestartSec=2
StartLimitBurst=3

# ログ設定
StandardOutput=journal
StandardError=journal

# プロセス管理
KillMode=control-group
TimeoutStopSec=10

# リソース制限（Raspberry Pi Zero 2 W向け）
MemoryMax=256M
CPUQuota=80%

[Install]
WantedBy=multi-user.target